# Memory-Monitor (内存监控) (Shell)/memory_monitor.timer
[Unit]
Description=Run memory monitor every minute

[Timer]
OnCalendar=*:0/1
Unit=memory_monitor.service

[Install]
WantedBy=timers.target
